# ops/strangelab-control-plane/services/systemd/strangelab-god-gateway.service
[Unit]
Description=StrangeLab God Button HTTP Gateway
After=network-online.target strangelab-token.service
Wants=network-online.target

[Service]
WorkingDirectory=/opt/strangelab
ExecStart=/usr/bin/node /opt/strangelab/god-gateway.mjs
Restart=always
RestartSec=2
Environment=NODE_ENV=production
Environment=LIVEKIT_URL=ws://192.168.8.114:7880
Environment=TOKEN_ENDPOINT=http://192.168.8.114:9123/token
Environment=ROOM_NAME=strangelab
Environment=IDENTITY=god-gateway
Environment=PORT=8099
Environment=HOST=0.0.0.0

[Install]
WantedBy=multi-user.target
