# ops/strangelab-control-plane/services/systemd/strangelab-ops-feed.service
[Unit]
Description=StrangeLab Ops Feed API
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/strangelab
User=pi
Group=pi
ExecStart=/usr/bin/node /opt/strangelab/ops-feed.mjs
Restart=always
RestartSec=2
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=9101
Environment=READ_MODE=auto
Environment=OPS_FEED_READ_TIMEOUT_MS=6000
Environment=OPS_FEED_HEALTH_CACHE_MS=5000
Environment=OPS_FEED_SSH_CONNECT_TIMEOUT_S=5
Environment=OPS_FEED_SSH_RETRIES=1
Environment=REMOTE_HOST=pi@pi-logger.local
Environment=LOGGER_HOST=pi-logger.local
Environment=REMOTE_EVENTS_DIR=/vault/sods/vault/events
Environment=SL_SSH_BIN=/usr/local/bin/sl-ssh
Environment=SL_SSH_ALIAS=vault

[Install]
WantedBy=multi-user.target
